# contrib/init/animusd.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env ANSD_BIN="/usr/bin/animusd"
env ANSD_USER="animus"
env ANSD_GROUP="animus"
env ANSD_PIDDIR="/var/run/animusd"
# upstart can't handle variables constructed with other variables
env ANSD_PIDFILE="/var/run/animusd/animusd.pid"
env ANSD_CONFIGFILE="/etc/animus/animus.conf"
env ANSD_DATADIR="/var/lib/animusd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # animusd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$ANSD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run animusd."
        echo "The setting must appear in $ANSD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$ANSD_PIDDIR"
    chmod 0755 "$ANSD_PIDDIR"
    chown $ANSD_USER:$ANSD_GROUP "$ANSD_PIDDIR"
    chown $ANSD_USER:$ANSD_GROUP "$ANSD_CONFIGFILE"
    chmod 0660 "$ANSD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$ANSD_PIDFILE" \
    --chuid $ANSD_USER:$ANSD_GROUP \
    --exec "$ANSD_BIN" \
    -- \
    -pid="$ANSD_PIDFILE" \
    -conf="$ANSD_CONFIGFILE" \
    -datadir="$ANSD_DATADIR" \
    -disablewallet \
    -daemon

